# rootfs/etc/s6-overlay/s6-rc.d/init/up
#!/command/with-contenv sh
exec /docker/init.sh
